# deployment/scipy2014.conf
[program:scipy2014]
command = /home/scipy/site/bin/runserver.sh
user = scipy
stdout_logfile = /home/scipy/site/logs/gunicorn_supervisor.log
redirect_stderr = true
autostart = true
killasgroup = true
